# Stock-Analysis/src/Ksfraser/Finance/2000/OtherSoftware/ast/share/example/crontab
*/5 * * * * (cd /var/www/localhost/htdocs/ast ; ./ast --scheduler)
